# hosts/strelok/systemd/user/sync_backup.timer
[Unit]
Description=Daily sync of all database backups from jarvis

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
